# bot.service
[Unit]
Description=Telegram Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/BOTtgOlegS
Environment="PATH=/usr/bin:/usr/local/bin"
# Используем python из venv (если venv создан)
ExecStart=/root/BOTtgOlegS/venv/bin/python3 /root/BOTtgOlegS/bot.py
# Если venv не используется, раскомментируйте следующую строку и закомментируйте выше:
# ExecStart=/usr/bin/python3 /root/BOTtgOlegS/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target